# dot_config/profile.d/20-emacs-bin.sh
appendpath /home/kvineet/.config/emacs/bin
