# dot_config/profile.d/30-krew.sh.tmpl
{{- if index . "krewPackages" }}#!/usr/bin/env sh
export KREW_ROOT="{{ .xdg.data }}/krew"
export PATH="$KREW_ROOT/bin:$PATH"
{{ end }}
