# dot_config/profile.d/30-krew.sh
#!/usr/bin/env sh
export KREW_ROOT="/home/kvineet/.local/share/krew"
export PATH="${KREW_ROOT:-$HOME/.krew}/bin:$PATH"
