# dot_config/profile.d/40-nix.sh
#!/usr/bin/env sh


if [ -e /home/kvineet/.nix-profile/etc/profile.d/nix.sh ]; then . /home/kvineet/.nix-profile/etc/profile.d/nix.sh; fi
source $HOME/.nix-profile/etc/profile.d/nix.sh
