# dot_config/profile.d/50-nvm.sh.tmpl
#!/usr/bin/env sh

export NVM_DIR="{{ .xdg.data }}/nvm"
[ -s "$NVM_DIR/nvm.sh" ] && \. "$NVM_DIR/nvm.sh"  # This loads nvm
